# src/rootfs/etc/services.d/inverter/finish
#!/usr/bin/with-contenv bashio
# finish скрипт для s6-overlay

exec 2>&1
bashio::log.info "Stopping Voltronic addon..."